# systemd/milvus-rag.service
[Unit]
Description=Milvus RAG API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Ortam değişkenlerini buradan oku (MILVUS_RAG_ROOT, MILVUS_RAG_LOG_DIR, MILVUS_HOST, MILVUS_PORT)
EnvironmentFile=/etc/default/milvus-rag

# Login shell ile çalıştırıp dizine geç, ortamdan belirlenen Python ile başlat
ExecStart=/bin/bash -lc 'cd "$MILVUS_RAG_ROOT" && exec "${MILVUS_RAG_PYTHON:-/usr/bin/python3}" -m gunicorn --bind 0.0.0.0:5000 --workers 4 --timeout 120 --access-logfile "$MILVUS_RAG_LOG_DIR"/access.log --error-logfile "$MILVUS_RAG_LOG_DIR"/error.log app:app'
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target