# systemd/milvus-rag.service
[Unit]
Description=Milvus RAG API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/opt/milvus-rag
Environment=PATH=/opt/milvus-rag/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/milvus-rag
Environment=MILVUS_HOST=localhost
Environment=MILVUS_PORT=19530
ExecStart=/opt/milvus-rag/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 4 --timeout 120 --access-logfile /opt/milvus-rag/logs/access.log --error-logfile /opt/milvus-rag/logs/error.log app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/milvus-rag/logs

[Install]
WantedBy=multi-user.target